# configurations/applications/zsh/scripts/profile.d/variables.sh
#default variables
#JAVA_HOME="/usr/lib/jvm/java-8-oracle/jre/bin/java"
#PATH="$PATH":/opt/idea-IC-145.1617.8/bin
export EDITOR=vim
